# systemd/battery-alert.service
[Unit]
Description=Battery Low Alert Service

[Service]
Type=oneshot
ExecStart=/home/exul4nzs/.battery_alert.sh
